# etc/systemd/user/openrgb.service
[Unit]
Description=OpenRGB

[Service]
ExecStart=/usr/bin/openrgb --server

[Install]
WantedBy=default.target
